# infra/ansible/roles/fortgym/templates/fort-gym-api.service.j2
[Unit]
Description=fort-gym API server
After=network.target dfhack-headless.service

[Service]
EnvironmentFile={{ fortgym_env_path }}
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ fortgym_install_dir }}
ExecStart={{ fortgym_venv_dir }}/bin/uvicorn fort_gym.bench.api.server:app --host 127.0.0.1 --port {{ fortgym_service_port }} --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
